# scripts/release-expired-tickets.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * /scripts/release-expired-tickets.sh
